# dashcam/board/raspberrypicm4io_64_clean/overlays/common/etc/systemd/system/generate-ssid.service
[Unit]
Description=Generate Unique Wi-Fi SSID from Device Serial
Before=hostapd.service
After=systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=/opt/dashcam/bin/generate-ssid.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target